# systemd/iphone-backup-status.service
[Unit]
Description=iPhone Backup Status Dashboard
After=network.target

[Service]
ExecStart=/usr/bin/authbind --deep /usr/bin/python3 /opt/iphone-backup/status-server.py
Restart=always
RestartSec=5
User=STATUS_USER
Environment=BACKUP_LOG=/var/log/iphone-backup.log
Environment=BORG_REPO=/backups/borg
Environment=STATUS_PORT=443

[Install]
WantedBy=multi-user.target
